# mikrotik-2fa-vpn.service
[Unit]
Description=MikroTik 2FA VPN System
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/mikrotik-2fa-vpn
Environment="PATH=/root/mikrotik-2fa-vpn/venv/bin"
ExecStart=/root/mikrotik-2fa-vpn/venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mikrotik-2fa-vpn

# Безопасность
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
