# mikrotik-2fa-vpn-telegram-bot.service
[Unit]
Description=MikroTik 2FA VPN Telegram Bot
After=network.target mikrotik-2fa-vpn.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/mikrotik-2fa-vpn
Environment="PATH=/root/mikrotik-2fa-vpn/venv/bin"
ExecStart=/root/mikrotik-2fa-vpn/venv/bin/python3 -m telegram_bot.bot
Restart=always
RestartSec=10

# Логирование
StandardOutput=append:/root/mikrotik-2fa-vpn/logs/telegram_bot.log
StandardError=append:/root/mikrotik-2fa-vpn/logs/telegram_bot.log
SyslogIdentifier=mikrotik-2fa-vpn-telegram-bot

# Безопасность
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
